# packaging/deb/postrm
#!/bin/sh


if ["$1" = "remove"]; then
    update-rc.d -f fdb-sql-layer remove
fi

